# apps/macos/packaging/com.claudex.claudexd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  LaunchAgent template for the Claudex engine-service (claudexd). Install per-user at
  ~/Library/LaunchAgents/com.claudex.claudexd.plist, then `launchctl load` it (or let the
  app offer to install it). The app itself does NOT need this — it will connect to claudexd
  whenever it is running — but a LaunchAgent keeps the durable engine alive across logins.

  Substitute the __PLACEHOLDERS__ for your install (the build script / installer does this):
    __NODE_BIN__   absolute path to the notarized Node (e.g. ~/.claudex/node/bin/node)
    __CLAUDEXD_JS__ absolute path to packages/cli/dist/claudexd.js (or a bundled copy)
    __HOME__       the user's home directory

  Modern alternative: SMAppService.agent(plistName:) registers this from inside the app
  bundle without writing to ~/Library/LaunchAgents directly. See README.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claudex.claudexd</string>
    <key>ProgramArguments</key>
    <array>
        <string>__NODE_BIN__</string>
        <string>__CLAUDEXD_JS__</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>ProcessType</key>
    <string>Background</string>
    <key>StandardOutPath</key>
    <string>__HOME__/.claudex/daemon/claudexd.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/.claudex/daemon/claudexd.log</string>
</dict>
</plist>
